# buildroot_fs_overlay/etc/profile.d/qt_evdev.sh
# see https://doc.qt.io/qt-5/inputs-linux-device.html
export QT_QPA_EVDEV_KEYBOARD_PARAMETERS=/dev/input/event0
export QT_QPA_EVDEV_MOUSE_PARAMETERS=/dev/input/event1:abs
